# systemd/user/kitty-theme@.service
[Unit]
Description=Set Kitty terminal theme to %i

[Service]
Type=oneshot
ExecStart=/home/himanshu/.config/kitty/toggle-theme.sh %i
